# pipeline/scripts/systemd/anka-ta-karpathy-open.timer
[Unit]
Description=Run Anka TA Karpathy holdout OPEN daily 09:15 IST

[Timer]
OnCalendar=Mon..Fri *-*-* 09:15:00
Persistent=true
Unit=anka-ta-karpathy-open.service

[Install]
WantedBy=timers.target
