# pipeline/scripts/systemd/anka-ta-karpathy-predict.timer
[Unit]
Description=Run Anka TA Karpathy predict daily 04:30 IST

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true
Unit=anka-ta-karpathy-predict.service

[Install]
WantedBy=timers.target
