# pipeline/backups/scheduled_tasks/2026-04-16/dusmtask.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.6" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <SecurityDescriptor>D:(A;;FA;;;SY)(A;;FA;;;BA)(A;;FA;;;S-1-5-80-4155767994-3874329934-3800885181-2130851812-726865888)(A;;FRFX;;;AU)(A;;FRFX;;;BU)(A;;FRFX;;;IU)</SecurityDescriptor>
    <URI>\Microsoft\Windows\DUSM\dusmtask</URI>
  </RegistrationInfo>
  <Principals>
    <Principal id="Users">
      <GroupId>S-1-5-4</GroupId>
    </Principal>
  </Principals>
  <Settings>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <MultipleInstancesPolicy>Queue</MultipleInstancesPolicy>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>
  </Settings>
  <Triggers>
    <WnfStateChangeTrigger>
      <StateName>7508BCA3291B950C</StateName>
    </WnfStateChangeTrigger>
  </Triggers>
  <Actions Context="Users">
    <Exec>
      <Command>%SystemRoot%\System32\dusmtask.exe</Command>
    </Exec>
  </Actions>
</Task>