# pipeline/backups/scheduled_tasks/2026-04-16/GatherNetworkInfo.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.6" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Source>$(@%SystemRoot%\system32\nettrace.dll,-6910)</Source>
    <Author>$(@%SystemRoot%\system32\nettrace.dll,-6911)</Author>
    <Description>$(@%SystemRoot%\system32\nettrace.dll,-6912)</Description>
    <URI>\Microsoft\Windows\NetTrace\GatherNetworkInfo</URI>
  </RegistrationInfo>
  <Principals>
    <Principal id="AnyUser">
      <GroupId>S-1-5-32-545</GroupId>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>true</StopIfGoingOnBatteries>
    <MultipleInstancesPolicy>Parallel</MultipleInstancesPolicy>
    <IdleSettings>
      <StopOnIdleEnd>true</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>
  </Settings>
  <Triggers />
  <Actions Context="AnyUser">
    <Exec>
      <Command>%windir%\system32\gatherNetworkInfo.vbs</Command>
      <WorkingDirectory>$(Arg1)</WorkingDirectory>
    </Exec>
  </Actions>
</Task>