# pipeline/backups/scheduled_tasks/2026-04-16/MareBackup.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Version>1.0</Version>
    <SecurityDescriptor>D:(A;;GA;;;BA)(A;;GA;;;SY)(A;;GA;;;BU)(A;;FRFX;;;LS)</SecurityDescriptor>
    <Source>$(@%SystemRoot%\system32\invagent.dll,-701)</Source>
    <Author>$(@%SystemRoot%\system32\invagent.dll,-701)</Author>
    <Description>$(@%SystemRoot%\system32\invagent.dll,-702)</Description>
    <URI>\Microsoft\Windows\Application Experience\MareBackup</URI>
  </RegistrationInfo>
  <Principals>
    <Principal id="LocalSystem">
      <UserId>S-1-5-18</UserId>
    </Principal>
  </Principals>
  <Settings>
    <DisallowStartIfOnBatteries>true</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>true</StopIfGoingOnBatteries>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <IdleSettings>
      <StopOnIdleEnd>true</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>
  </Settings>
  <Triggers />
  <Actions Context="LocalSystem">
    <Exec>
      <Command>%windir%\system32\compattelrunner.exe</Command>
      <Arguments>-m:aeinv.dll -f:UpdateSoftwareInventoryW invsvc</Arguments>
    </Exec>
    <Exec>
      <Command>%windir%\system32\compattelrunner.exe</Command>
      <Arguments>-m:appraiser.dll -f:DoScheduledTelemetryRun</Arguments>
    </Exec>
    <Exec>
      <Command>%windir%\system32\compattelrunner.exe</Command>
      <Arguments>-m:aemarebackup.dll -f:BackupMareData</Arguments>
    </Exec>
  </Actions>
</Task>