# pipeline/backups/scheduled_tasks/2026-04-16/UPnPHostConfig.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.6" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <SecurityDescriptor>D:(A;;FA;;;BA)(A;;FA;;;SY)(A;;FRFX;;;LS)</SecurityDescriptor>
    <Author>$(@%systemroot%\system32\upnphost.dll,-215)</Author>
    <Description>$(@%systemroot%\system32\upnphost.dll,-216)</Description>
    <URI>\Microsoft\Windows\UPnP\UPnPHostConfig</URI>
  </RegistrationInfo>
  <Principals>
    <Principal id="LocalSystem">
      <UserId>S-1-5-18</UserId>
    </Principal>
  </Principals>
  <Settings>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <IdleSettings>
      <StopOnIdleEnd>true</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>
  </Settings>
  <Triggers />
  <Actions Context="LocalSystem">
    <Exec>
      <Command>sc.exe</Command>
      <Arguments>config upnphost start= auto</Arguments>
    </Exec>
  </Actions>
</Task>