# pipeline/infra/systemd/anka-security-daily.timer
[Unit]
Description=Anka — Daily security timer (06:00 IST)

[Timer]
OnCalendar=*-*-* 06:00:00
AccuracySec=15min
Unit=anka-security-daily.service
Persistent=true

[Install]
WantedBy=timers.target
