# pipeline/infra/systemd/anka-gemma4-auto-disable.service
[Unit]
Description=Anka Gemma 4 Pilot — hourly auto-disable guardrail (09:00–22:00 IST)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=anka
Group=anka
WorkingDirectory=/home/anka/askanka.com
EnvironmentFile=/home/anka/askanka.com/pipeline/.env
ExecStart=/home/anka/askanka.com/.venv/bin/python pipeline/scripts/gemma4_auto_disable_check.py
StandardOutput=append:/home/anka/askanka.com/opus/logs/gemma4_pilot.log
StandardError=append:/home/anka/askanka.com/opus/logs/gemma4_pilot.log
Nice=10
TimeoutStartSec=5min
