# pipeline/infra/systemd/anka-gemma4-daily-report.timer
[Unit]
Description=Run Anka Gemma 4 daily report at 22:00 IST

[Timer]
OnCalendar=*-*-* 22:00:00 Asia/Kolkata
Persistent=true
Unit=anka-gemma4-daily-report.service

[Install]
WantedBy=timers.target
