# pipeline/infra/systemd/anka-secrsi-basket-close.service
[Unit]
Description=Anka SECRSI — 14:30 IST mechanical TIME_STOP close
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=anka
Group=anka
WorkingDirectory=/home/anka/askanka.com
EnvironmentFile=/home/anka/askanka.com/pipeline/.env
Environment=PYTHONPATH=/home/anka/askanka.com
ExecStart=/home/anka/askanka.com/.venv/bin/python -m pipeline.research.h_2026_04_27_secrsi.forward_shadow basket-close
StandardOutput=append:/home/anka/askanka.com/pipeline/logs/secrsi.log
StandardError=append:/home/anka/askanka.com/pipeline/logs/secrsi.log
Nice=10
TimeoutStartSec=10min
