# pipeline/infra/systemd/anka-auto-push.service
[Unit]
Description=Anka — Auto-push all local branches to GitHub
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=anka
Group=anka
Environment=ANKA_REPO_ROOT=/home/anka/askanka.com
ExecStart=/home/anka/askanka.com/pipeline/scripts/auto_push_branches.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
